# deploy/systemd/connwatch-agent.service
[Unit]
Description=Connwatch Agent
After=network-online.target
Wants=network-online.target

[Service]
User=connwatch
Group=connwatch
ExecStart=/usr/local/bin/connwatch-agent /etc/connwatch-agent.yaml
Restart=always
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
SystemCallFilter=@system-service
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
